# couchdb/couchdb.ini
[program:couchdb]
command=couchdb
user=couchdb
