# kubernetes/cfg_files/supervisord.d/kubernetes.ini
[program:kubernetes]
command=/usr/local/bin/runKube
